# scripts/sopds.service
[Unit]
Description=SOPDS NG OPDS catalog server
Documentation=https://github.com/sarutobi/sopds-ng
After=network.target
Wants=network.target

[Service]
Type=simple
User=sopds
Group=sopds
WorkingDirectory=/opt/sopds-ng

EnvironmentFile=-/data/.env
Environment=DATA_ROOT=/data
Environment=PYTHONDONTWRITEBYTECODE=1

ExecStart=/opt/sopds-ng/scripts/start_server.sh

Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/data

LimitNOFILE=65535
LimitNPROC=4096

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sopds

[Install]
WantedBy=multi-user.target
